# src/ip-update/ip-update.timer
[Unit]
Description=IP Update

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
